# root/etc/s6-overlay/s6-rc.d/init-config-ml/run
#!/usr/bin/with-contenv bash

mkdir -p /config/machine-learning

# apply early ml support
cp /defaults/main.py /app/immich/machine-learning/src/main.py
if [ ! -d /config/cuda-venv ]; then
    python3 -m venv /config/cuda-venv/
fi

# install/update pip packages to /config venv
source /config/cuda-venv/bin/activate

pip install -U --no-cache-dir --index-url https://download.pytorch.org/whl/cu117 \
    torch==${TORCH_VERSION}

pip install -U --no-cache-dir \
    fastapi \
    insightface \
    nltk \
    numpy \
    onnxruntime-gpu \
    packaging \
    Pillow \
    scikit-learn \
    scipy \
    sentencepiece \
    sentence-transformers \
    tqdm \
    transformers \
    uvicorn[standard]

# download and install cudnn
if [ ! -f /config/lib/libcudnn8_8.5.0.96-1+cuda11.7_amd64.deb ]; then
    rm -f /config/lib/libcudnn*.deb
    mkdir -p /config/lib
    curl -o \
        /config/lib/libcudnn8_8.5.0.96-1+cuda11.7_amd64.deb -L \
        https://developer.download.nvidia.com/compute/cuda/repos/ubuntu2204/x86_64/libcudnn8_8.5.0.96-1+cuda11.7_amd64.deb
fi
dpkg -i /config/lib/libcudnn8_8.5.0.96-1+cuda11.7_amd64.deb

# permissions
chown -R abc:abc \
    /config \
    /app
