# root/etc/s6-overlay/s6-rc.d/init-config-ml/run
#!/usr/bin/with-contenv bash

mkdir -p /config/machine-learning

# apply early ml support
cp /defaults/main.py /app/immich/machine-learning/main.py
if [ ! -d /config/cuda-venv ]; then
    python3 -m venv /config/cuda-venv/
fi

# install/update pip packages to /config venv
source /config/cuda-venv/bin/activate

pip install -U --no-cache-dir \
    aiocache \
    fastapi \
    insightface \
    nltk \
    numpy \
    onnxruntime-gpu \
    Pillow \
    scikit-learn \
    scipy \
    sentencepiece \
    sentence-transformers \
    toch \
    tqdm \
    transformers \
    uvicorn[standard]

# permissions
chown -R abc:abc \
    /photos \
    /config \
    /app
